# overlay/etc/cron.hourly/magento
#!/bin/bash

/usr/bin/php /var/www/magento/cron.php > /dev/null 2>&1
